# unit_client.service
[Unit]
Description=Sys_Logger Unit Client
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /path/to/unit_client.py --start
ExecStop=/usr/bin/python3 /path/to/unit_client.py --stop
Restart=always
RestartSec=10
Environment=PYTHONPATH=/path/to
WorkingDirectory=/path/to

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/tmp /var/tmp /run

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unit_client

[Install]
WantedBy=multi-user.target